# scripts/userscripts/set-sleep-timer.service
[Unit]
Description=Set sleep timer with volume reduction

[Service]
ExecStart=/home/pi/RPi-Jukebox-RFID/scripts/userscripts/set-sleep-timer.sh